# kiosk-stable.service
[Unit]
Description=Ultra-Stable Smart Display Kiosk Browser
After=network.target pi-smart-display.service
Wants=pi-smart-display.service

[Service]
Type=simple
ExecStart=/home/ross/pi-smart-display/stable-kiosk.sh
WorkingDirectory=/home/ross/pi-smart-display
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5
User=ross
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Install]
WantedBy=graphical.target
