# motion.service
[Unit]
Description=PIR Motion Detection for Smart Display
After=network.target

[Service]
Type=simple
User=ross
WorkingDirectory=/home/ross/pi-smart-display
ExecStart=/usr/bin/python3 /home/ross/pi-smart-display/motion-detection.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
