# resolution.service
[Unit]
Description=Set 4K Resolution on Boot
After=graphical.target

[Service]
Type=oneshot
ExecStart=/home/ross/pi-smart-display/set-4k-resolution.sh
User=ross
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
